# systemd/tunnelfw-sync.timer
[Unit]
Description=Periodic tunnelfw AD membership sync

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
